# ansible/roles/db_backup/templates/git-db-backup.timer.j2
[Unit]
Description=Daily git-push job for git-db-backup.service

[Timer]
# daily
OnCalendar=*-*-* 17:30:00
# hourly
# OnCalendar=*-*-* *:00:00
Persistent=false


[Install]
WantedBy=multi-user.target
